# dotfiles/immutable/desktop/.config/autostart/net.opentabletdriver.OpenTabletDriver.desktop
[Desktop Entry]
Type=Application
Name=OpenTabletDriver
Exec=flatpak run net.opentabletdriver.OpenTabletDriver -d
