# dotfiles/.config/autostart/kdeconnect-indicator.desktop
[Desktop Entry]
Type=Application
Name=kdeconnect-indicator
Exec=kdeconnect-indicator
