# dotfiles/.config/autostart/shepherd.desktop
# Start user Shepherd in environments that support XDG Autostart.
[Desktop Entry]
Type=Application
Name=Shepherd
Exec=shepherd
